# supervisor.conf
[supervisord]
nodaemon=true

[program:mongodb]
user=root
command=/usr/bin/mongod
redirect_stderr=true

[program:chat]
user=rocketchat
command=/app/startup.sh
startsecs=5
redirect_stderr=true

